# playbooks/roles/s3-backup/templates/s3-backup.timer.j2
[Unit]
Description=Backup folders to s3

[Timer]
OnCalendar=daily
Unit=s3-backup.service

[Install]
WantedBy=timers.target
