# playbooks/roles/s3-backup/templates/s3-backup.service.j2
[Unit]
Description=Backup folders to s3
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/s3-backup.sh
RemainAfterExit=yes
